# configs/openfaas/startup-prometheus.conf
start on startup
task
/etc/init/startup-prometheus.sh
