# configs/openfaas/startup-prometheus.conf
start on startup
task
~/startup-prometheus.sh
