# configs/openfaas/startup-prometheus.conf
start on startup
task
apt-get update
apt-get install wget
/etc/init/startup-prometheus.sh
